# root/etc/hotplug.d/iface/97-parentcontrol
#!/bin/sh
CONFIG=parentcontrol

[ "$ACTION" = "ifup" ] || exit 0
[ "$(uci -q get $CONFIG.@basic[0].enabled)" = "1" ] || exit 0

# Use reload (lightweight) instead of start to avoid rebuilding cron on every ifup
/etc/init.d/parentcontrol reload