# root/etc/init.d/parentcontrol
#!/bin/sh /etc/rc.common
# ParentControl - nftables fw4 version (cron for time control)

START=98
STOP=15
CONFIG=parentcontrol
LOG_DIR=/tmp/parentcontrol
LOG_FILE=$LOG_DIR/$CONFIG.log
LOCK=/tmp/lock/${CONFIG}.lock
STATE_FILE=$LOG_DIR/rules_state
nft=$(command -v nft)
[ -z "$nft" ] && nft="/usr/sbin/nft"
NOW_MIN=0
NOW_WEEK=0
NOW_H=0
NOW_M=0

EXTRA_COMMANDS="status"
EXTRA_HELP="        status  Show service status"

elog() {
    mkdir -p "$LOG_DIR"
    local d="$(date "+%Y-%m-%d %H:%M:%S")"
    printf "%s: %s\n" "$d" "$*" >>"$LOG_FILE"
}

clean_log() {
    [ -f "$LOG_FILE" ] && [ "$(wc -l < "$LOG_FILE")" -gt 500 ] && {
        tail -100 "$LOG_FILE" > "$LOG_FILE.tmp"
        mv "$LOG_FILE.tmp" "$LOG_FILE"
        elog "log trimmed"
    }
}

in_week() {
    local week_str="$1" day="$2"
    [ -z "$week_str" ] && return 0
    [ "$week_str" = "*" ] && return 0
    for w in ${week_str//,/ }; do
        [ "$w" = "$day" ] && return 0
    done
    return 1
}

time_match() {
    local start="$1" end="$2" week="$3" now_min="$4" now_week="$5"
    in_week "$week" "$now_week" || return 1
    [ -z "$start" ] && [ -z "$end" ] && return 0
    [ -z "$start" ] || [ -z "$end" ] && return 0

    local sh sm eh em s e
    sh=${start%:*}; sm=${start#*:}; eh=${end%:*}; em=${end#*:}
    [ -z "$sh" ] || [ -z "$sm" ] || [ -z "$eh" ] || [ -z "$em" ] && return 1
    sh=${sh#0}; sm=${sm#0}; eh=${eh#0}; em=${em#0}
    [ -z "$sh" ] && sh=0; [ -z "$sm" ] && sm=0; [ -z "$eh" ] && eh=0; [ -z "$em" ] && em=0
    s=$((sh * 60 + sm))
    e=$((eh * 60 + em))

    if [ "$e" -ge "$s" ]; then
        [ "$now_min" -ge "$s" ] && [ "$now_min" -le "$e" ]
        return $?
    else
        # Cross-midnight range, e.g. 22:00-08:30
        [ "$now_min" -ge "$s" ] || [ "$now_min" -le "$e" ]
        return $?
    fi
}

update_now() {
    local h m
    h=$(date +%H); m=$(date +%M)
    h=${h#0}; m=${m#0}
    [ -z "$h" ] && h=0; [ -z "$m" ] && m=0
    NOW_H=$h; NOW_M=$m
    NOW_MIN=$((h * 60 + m))
    NOW_WEEK=$(date +%w)
    [ "$NOW_WEEK" = "0" ] && NOW_WEEK=7
}

flush_conntrack() {
    command -v conntrack >/dev/null 2>&1 && conntrack -F 2>/dev/null
    local ft
    for ft in $($nft list flowtables inet fw4 2>/dev/null | awk '/flowtable /{print $2}'); do
        $nft flush flowtable inet fw4 "$ft" 2>/dev/null
    done
}

remove_jump_rules() {
    local target_chain="$1" base_chain="${2:-forward}"
    local h
    for h in $($nft -a list chain inet fw4 "$base_chain" 2>/dev/null \
                | grep "jump $target_chain" | sed 's/.*# handle //'); do
        $nft delete rule inet fw4 "$base_chain" handle "$h" 2>/dev/null
    done
}

setup_chains() {
    local control_mode
    control_mode=$(uci -q get $CONFIG.@basic[0].control_mode)
    [ -z "$control_mode" ] && control_mode=0

    $nft add table inet fw4 2>/dev/null || true
    $nft add chain inet fw4 parentcontrol_time 2>/dev/null || true
    $nft add chain inet fw4 parentcontrol_speed 2>/dev/null || true

    local ch
    for ch in parentcontrol_time parentcontrol_speed; do
        remove_jump_rules "$ch" forward
        $nft insert rule inet fw4 forward jump "$ch" 2>/dev/null
    done

    remove_jump_rules parentcontrol_time input
    if [ "$control_mode" = "1" ]; then
        $nft insert rule inet fw4 input jump parentcontrol_time 2>/dev/null
    fi
}

del_rule() {
    local ch
    for ch in parentcontrol_time parentcontrol_speed; do
        remove_jump_rules "$ch" forward
        remove_jump_rules "$ch" input
    done
    for ch in parentcontrol_time parentcontrol_speed; do
        $nft flush chain inet fw4 "$ch" 2>/dev/null
        $nft delete chain inet fw4 "$ch" 2>/dev/null
    done
}

# Compute a string representing which rules should be active now
# Used to detect state changes and avoid unnecessary conntrack flushes
compute_rules_state() {
    local state="" idx mmode
    for mmode in time speed; do
        idx=0
        while uci -q get "$CONFIG.@${mmode}[$idx]" >/dev/null 2>&1; do
            local enabled tstart tstop week mac
            enabled=$(uci -q get "$CONFIG.@${mmode}[$idx].enable")
            if [ "$enabled" = "1" ]; then
                tstart=$(uci -q get "$CONFIG.@${mmode}[$idx].timestart")
                tstop=$(uci -q get "$CONFIG.@${mmode}[$idx].timeend")
                week=$(uci -q get "$CONFIG.@${mmode}[$idx].week")
                mac=$(uci -q get "$CONFIG.@${mmode}[$idx].mac")

                local active=1
                if [ -n "$tstart" ] && [ -n "$tstop" ]; then
                    time_match "$tstart" "$tstop" "$week" "$NOW_MIN" "$NOW_WEEK" || active=0
                fi

                if [ "$active" = "1" ]; then
                    case "$mmode" in
                        time)  state="${state}T${idx}:${mac:-ALL};" ;;
                        speed)
                            local rate
                            rate=$(uci -q get "$CONFIG.@${mmode}[$idx].rate")
                            state="${state}S${idx}:${mac:-ALL}:${rate:-1mbit};"
                        ;;
                    esac
                fi
            fi
            idx=$((idx + 1))
        done
    done
    echo "$state"
}

# Load nft rules for all currently-active entries
load_active_rules() {
    local mmode="$1" idx=0
    local lan_if
    if [ "$mmode" = "time" ]; then
        lan_if=$(uci -q get network.lan.device)
        [ -z "$lan_if" ] && lan_if=$(uci -q get network.lan.ifname)
        [ -z "$lan_if" ] && lan_if="br-lan"
    fi

    while uci -q get "$CONFIG.@${mmode}[$idx]" >/dev/null 2>&1; do
        local enabled tstart tstop week mac
        enabled=$(uci -q get "$CONFIG.@${mmode}[$idx].enable")
        if [ "$enabled" = "1" ]; then
            tstart=$(uci -q get "$CONFIG.@${mmode}[$idx].timestart")
            tstop=$(uci -q get "$CONFIG.@${mmode}[$idx].timeend")
            week=$(uci -q get "$CONFIG.@${mmode}[$idx].week")

            local active=1
            if [ -n "$tstart" ] && [ -n "$tstop" ]; then
                time_match "$tstart" "$tstop" "$week" "$NOW_MIN" "$NOW_WEEK" || active=0
            fi

            if [ "$active" = "1" ]; then
                mac=$(uci -q get "$CONFIG.@${mmode}[$idx].mac")
                local MAC_EXPR=""
                [ -n "$mac" ] && MAC_EXPR="ether saddr $mac"

                case "$mmode" in
                    time)
                        local INF_EXPR="meta iifname \"$lan_if\""
                        $nft add rule inet fw4 parentcontrol_time $INF_EXPR $MAC_EXPR counter reject 2>/dev/null
                        elog "time+ mac=${mac:-ALL} [$tstart-$tstop w=${week:-*}]"
                    ;;
                    speed)
                        local rate
                        rate=$(uci -q get "$CONFIG.@${mmode}[$idx].rate")
                        [ -z "$rate" ] && rate="1mbit"
                        $nft add rule inet fw4 parentcontrol_speed $MAC_EXPR limit rate over $rate counter drop 2>/dev/null
                        elog "speed+ mac=${mac:-ALL} rate=$rate [$tstart-$tstop w=${week:-*}]"
                    ;;
                esac
            fi
        fi
        idx=$((idx + 1))
    done
}

# Single cron entry: reload every minute (reload is smart, skips if no change)
rebuild_cron() {
    sed -i '/# parentcontrol$/d' /etc/crontabs/root 2>/dev/null
    echo "* * * * * /etc/init.d/parentcontrol reload # parentcontrol" >> /etc/crontabs/root
    /etc/init.d/cron restart
}

start() {
    local enabled
    enabled=$(uci -q get $CONFIG.@basic[0].enabled)
    [ "$enabled" = "1" ] || { elog "start: disabled"; exit 1; }

    mkdir -p "$LOG_DIR"
    touch "$LOCK"
    elog "=== service start ==="

    del_rule
    setup_chains
    update_now

    load_active_rules time
    load_active_rules speed
    flush_conntrack

    # Save initial state
    compute_rules_state > "$STATE_FILE"
    elog "started (now=$NOW_H:$NOW_M week=$NOW_WEEK)"

    rebuild_cron
    clean_log
}

stop() {
    del_rule
    rm -f "$LOCK" "$STATE_FILE" 2>/dev/null
    sed -i '/# parentcontrol$/d' /etc/crontabs/root 2>/dev/null
    /etc/init.d/cron restart
    elog "=== service stopped ==="
}

reload() {
    local enabled
    enabled=$(uci -q get $CONFIG.@basic[0].enabled)
    if [ "$enabled" != "1" ]; then
        stop
        return
    fi

    mkdir -p "$LOG_DIR"
    update_now

    local new_state old_state
    new_state=$(compute_rules_state)
    old_state=$(cat "$STATE_FILE" 2>/dev/null)

    # Fast path: no change → do nothing
    if [ "$new_state" = "$old_state" ]; then
        return 0
    fi

    # State changed: reload rules
    setup_chains
    $nft flush chain inet fw4 parentcontrol_time 2>/dev/null
    $nft flush chain inet fw4 parentcontrol_speed 2>/dev/null

    load_active_rules time
    load_active_rules speed
    flush_conntrack

    echo "$new_state" > "$STATE_FILE"
    touch "$LOCK"
    elog "reload: rules changed (now=$NOW_H:$NOW_M week=$NOW_WEEK) [$new_state]"
}

status() {
    local enabled
    enabled=$(uci -q get $CONFIG.@basic[0].enabled)

    if [ "$enabled" != "1" ]; then
        echo "parentcontrol is inactive (disabled)"
        return 1
    fi

    if [ ! -f "$LOCK" ]; then
        echo "parentcontrol is inactive (not started)"
        return 1
    fi

    local time_rules=0 speed_rules=0
    if $nft list chain inet fw4 parentcontrol_time >/dev/null 2>&1; then
        time_rules=$($nft list chain inet fw4 parentcontrol_time 2>/dev/null | grep -c "reject")
    fi
    if $nft list chain inet fw4 parentcontrol_speed >/dev/null 2>&1; then
        speed_rules=$($nft list chain inet fw4 parentcontrol_speed 2>/dev/null | grep -c "drop")
    fi

    if [ "$time_rules" -gt 0 ] || [ "$speed_rules" -gt 0 ]; then
        echo "parentcontrol is running (time=$time_rules speed=$speed_rules)"
    elif grep -q "# parentcontrol" /etc/crontabs/root 2>/dev/null; then
        echo "parentcontrol is running (cron active, waiting for time window)"
    else
        echo "parentcontrol is started (no active rules)"
    fi
    return 0
}