# root/etc/init.d/parentcontrol
#!/bin/sh /etc/rc.common
# ParentControl - nftables fw4 version (cron for time control)

START=98
CONFIG=parentcontrol
LOG_DIR=/tmp/parentcontrol
LOG_FILE=$LOG_DIR/$CONFIG.log
LOCK=/tmp/lock/${CONFIG}.lock
nft=$(command -v nft)
DNSMASQ_CONF_DIR=/tmp/dnsmasq.d
DNS_CONF=$DNSMASQ_CONF_DIR/parentcontrol.conf
TIME_FILTER=0
NOW_MIN=0
NOW_WEEK=0

EXTRA_COMMANDS="block_time unblock_time block_protocol unblock_protocol block_weburl unblock_weburl"

elog() {
    mkdir -p "$LOG_DIR"
    local d="$(date "+%Y-%m-%d %H:%M:%S")"
    printf "%s: %s\n" "$d" "$*" >>"$LOG_FILE"
}

clean_log() {
    [ -f "$LOG_FILE" ] && [ "$(wc -l < "$LOG_FILE")" -gt 500 ] && elog "clear log ."
}

in_week() {
    local week_str="$1" day="$2"
    [ -z "$week_str" ] && return 0
    [ "$week_str" = "*" ] && return 0
    for w in ${week_str//,/ }; do
        [ "$w" = "$day" ] && return 0
    done
    return 1
}

time_match() {
    # 参数：start end week now_min now_week
    local start="$1" end="$2" week="$3" now_min="$4" now_week="$5"

    # 星期不匹配
    in_week "$week" "$now_week" || return 1

    # 无时间段 -> 视为匹配
    [ -z "$start" ] && [ -z "$end" ] && return 0

    # 仅 start 或仅 end：不严控，视为匹配
    [ -z "$start" ] || [ -z "$end" ] && return 0

    local sh sm eh em s e
    sh=${start%:*}; sm=${start#*:}; eh=${end%:*}; em=${end#*:}
    # 基本校验，防止算术错误
    [ -z "$sh" ] || [ -z "$sm" ] || [ -z "$eh" ] || [ -z "$em" ] && return 1
    sh=${sh#0}; sm=${sm#0}; eh=${eh#0}; em=${em#0}
    [ -z "$sh" ] && sh=0; [ -z "$sm" ] && sm=0; [ -z "$eh" ] && eh=0; [ -z "$em" ] && em=0
    s=$((sh * 60 + sm))
    e=$((eh * 60 + em))

    if [ "$e" -ge "$s" ]; then
        # 同日区间
        [ "$now_min" -ge "$s" ] && [ "$now_min" -le "$e" ]
        return $?
    else
        # 跨午夜区间，例如 23:00-06:00
        [ "$now_min" -ge "$s" ] || [ "$now_min" -le "$e" ]
        return $?
    fi
}

# 生成/重建 cron：根据 time sections 写定时 block/unblock
rebuild_cron() {
    # 先移除旧的 parentcontrol 条目
    sed -i '/# parentcontrol$/d;/\/etc\/init\.d\/parentcontrol \(\<block_time\>\|\<unblock_time\>\|\<block_protocol\>\|\<unblock_protocol\>\|\<block_weburl\>\|\<unblock_weburl\>\)/d' /etc/crontabs/root 2>/dev/null

    cron_add() {
        local hhmm="$1"; local week="$2"; local cmd="$3"
        [ -z "$hhmm" ] && return
        echo "${hhmm#*:} ${hhmm%:*} * * $week $cmd # parentcontrol" >> /etc/crontabs/root
    }

    # 根据各 section 的 timestart/timeend 生成 cron
    for mode in time protocol weburl; do
        for i in $(uci show "$CONFIG" | grep "=${mode}" | cut -d'[' -f2 | cut -d']' -f1); do
            local enabled start end week
            enabled=$(uci -q get $CONFIG.@$mode[$i].enable)
            [ "x$enabled" != "x1" ] && continue

            start=$(uci -q get $CONFIG.@$mode[$i].timestart)
            end=$(uci -q get $CONFIG.@$mode[$i].timeend)
            week=$(uci -q get $CONFIG.@$mode[$i].week)
            [ -z "$week" ] && week="*"

            # 有时间段才写 cron
            if [ -n "$start" ] || [ -n "$end" ]; then
                cron_add "$start" "$week" "/etc/init.d/parentcontrol block_$mode"
                cron_add "$end"   "$week" "/etc/init.d/parentcontrol unblock_$mode"
            fi
        done
    done

    /etc/init.d/cron restart
}

set_rules() {
    mmode=$1
    timed_only=$2   # "1" 仅加载有时间段的条目；"0" 仅加载常驻；空/其他表示全部

    # 确保表和普通链存在
    $nft add table inet fw4 2>/dev/null || true
    case "$mmode" in
        time)     $nft add chain inet fw4 parentcontrol_time 2>/dev/null || true ;;
        protocol) $nft add chain inet fw4 parentcontrol_protocol 2>/dev/null || true ;;
        weburl)
            $nft add chain inet fw4 parentcontrol_weburl 2>/dev/null || true
            $nft add chain inet fw4 parentcontrol_weburl_nat 2>/dev/null || true
        ;;
    esac

    # 挂载 jump 到 base chain（显式插入到链首，避免被后续 ACCEPT 提前放行）
    case "$mmode" in
        time)
            $nft insert rule inet fw4 forward position 0 jump parentcontrol_time 2>/dev/null || \
                $nft add rule inet fw4 forward jump parentcontrol_time 2>/dev/null
        ;;
        protocol)
            $nft insert rule inet fw4 forward position 0 jump parentcontrol_protocol 2>/dev/null || \
                $nft add rule inet fw4 forward jump parentcontrol_protocol 2>/dev/null
        ;;
        weburl)
            $nft insert rule inet fw4 forward position 0 jump parentcontrol_weburl 2>/dev/null || \
                $nft add rule inet fw4 forward jump parentcontrol_weburl 2>/dev/null
            $nft insert rule inet fw4 prerouting position 0 jump parentcontrol_weburl_nat 2>/dev/null || \
                $nft add rule inet fw4 prerouting jump parentcontrol_weburl_nat 2>/dev/null
        ;;
    esac

    # 遍历配置，提取 section index
    for i in $(uci show "$CONFIG" | grep "=${mmode}" | cut -d'[' -f2 | cut -d']' -f1); do
        enabled=$(uci -q get $CONFIG.@$mmode[$i].enable)
        [ "x$enabled" != "x1" ] && continue

        start=$(uci -q get $CONFIG.@$mmode[$i].timestart)
        stop=$(uci -q get $CONFIG.@$mmode[$i].timeend)

        case "$timed_only" in
            1)  [ -z "$start" ] && [ -z "$stop" ] && continue ;;
            0)  [ -n "$start" ] || [ -n "$stop" ] || true; if [ -n "$start" ] || [ -n "$stop" ]; then continue; fi ;;
        esac

        # 需要按当前时间过滤的场景（定时条目 + TIME_FILTER=1 时）
        if [ "$timed_only" = "1" ] && [ "$TIME_FILTER" = "1" ]; then
            week=$(uci -q get $CONFIG.@$mmode[$i].week)
            time_match "$start" "$stop" "$week" "$NOW_MIN" "$NOW_WEEK" || continue
        fi

        mac=$(uci -q get $CONFIG.@$mmode[$i].mac)
        MAC_EXPR=""
        [ -n "$mac" ] && MAC_EXPR="ether saddr $mac"

        case "$mmode" in
            time)
                # 兼容 ImmortalWrt 24.10：优先使用 device，其次 ifname，缺省 br-lan
                lan_if=$(uci -q get network.lan.device)
                [ -z "$lan_if" ] && lan_if=$(uci -q get network.lan.ifname)
                [ -z "$lan_if" ] && lan_if="br-lan"
                INF_EXPR="meta iifname \"$lan_if\""
                # 仅阻断链，具体启停由 cron 驱动 block/unblock
                $nft add rule inet fw4 parentcontrol_time $INF_EXPR $MAC_EXPR counter reject 2>/dev/null
                elog "time-chain rule loaded: mac=$mac timed_only=$timed_only"
            ;;
            protocol)
                proto=$(uci -q get $CONFIG.@$mmode[$i].proto); [ -z "$proto" ] && proto="tcp"
                ports=$(uci -q get $CONFIG.@$mmode[$i].ports)
                dports=$(uci -q get $CONFIG.@$mmode[$i].portd)

                [ -n "$ports" ] && $nft add rule inet fw4 parentcontrol_protocol $MAC_EXPR $proto sport { $ports } counter reject 2>/dev/null
                [ -n "$dports" ] && $nft add rule inet fw4 parentcontrol_protocol $MAC_EXPR $proto dport { $dports } counter reject 2>/dev/null
                elog "protocol-chain rule loaded: mac=$mac proto=$proto sport=$ports dport=$dports timed_only=$timed_only"
            ;;
            weburl)
                # DNS 层方案：黑洞域名 + 强制 DNS 走路由器（含 DoT 阻断）
                lan_if=$(uci -q get network.lan.device)
                [ -z "$lan_if" ] && lan_if=$(uci -q get network.lan.ifname)
                [ -z "$lan_if" ] && lan_if="br-lan"
                router_ip=$(uci -q get network.lan.ipaddr)
                [ -z "$router_ip" ] && router_ip="192.168.1.1"

                # filter：丢弃 DoT；丢弃直连外部 DNS（IPv4/IPv6）
                $nft add rule inet fw4 parentcontrol_weburl $MAC_EXPR meta nfproto ipv4 udp dport 53 ip daddr != $router_ip counter drop 2>/dev/null
                $nft add rule inet fw4 parentcontrol_weburl $MAC_EXPR meta nfproto ipv6 udp dport 53 counter drop 2>/dev/null
                $nft add rule inet fw4 parentcontrol_weburl $MAC_EXPR tcp dport 853 counter drop 2>/dev/null

                # nat：把该 MAC 的 DNS 请求强制转发到路由器（仅 IPv4）
                $nft add rule inet fw4 parentcontrol_weburl_nat meta nfproto ipv4 iifname "$lan_if" $MAC_EXPR udp dport 53 dnat to $router_ip counter 2>/dev/null
                $nft add rule inet fw4 parentcontrol_weburl_nat meta nfproto ipv4 iifname "$lan_if" $MAC_EXPR tcp dport 53 dnat to $router_ip counter 2>/dev/null

                elog "weburl-chain rule loaded: mac=$mac (DNS redirect + DoT drop) timed_only=$timed_only"
            ;;
        esac
    done
}

build_dns_blacklist() {
    timed_only=$1

    mkdir -p "$DNSMASQ_CONF_DIR"
    : > "$DNS_CONF"

    idx=0
    for i in $(uci show "$CONFIG" | grep "=weburl" | cut -d'[' -f2 | cut -d']' -f1); do
        enabled=$(uci -q get $CONFIG.@weburl[$i].enable)
        [ "x$enabled" != "x1" ] && continue

        start=$(uci -q get $CONFIG.@weburl[$i].timestart)
        stop=$(uci -q get $CONFIG.@weburl[$i].timeend)
        week=$(uci -q get $CONFIG.@weburl[$i].week)
        case "$timed_only" in
            1)  [ -z "$start" ] && [ -z "$stop" ] && continue ;;
            0)  if [ -n "$start" ] || [ -n "$stop" ]; then continue; fi ;;
        esac

        if [ "$TIME_FILTER" = "1" ] && [ "$timed_only" = "1" ]; then
            time_match "$start" "$stop" "$week" "$NOW_MIN" "$NOW_WEEK" || continue
        fi

        mac=$(uci -q get $CONFIG.@weburl[$i].mac)
        ipaddr=$(uci -q get $CONFIG.@weburl[$i].ip)
        domain=$(uci -q get $CONFIG.@weburl[$i].word)
        [ -z "$mac" ] && continue
        [ -z "$domain" ] && continue

        tag="pcw_$idx"
        echo "dhcp-mac=set:$tag,$mac" >> "$DNS_CONF"
        [ -n "$ipaddr" ] && echo "dhcp-host=$mac,$ipaddr,set:$tag" >> "$DNS_CONF"
        echo "address=/$domain/0.0.0.0/tag:$tag" >> "$DNS_CONF"
        echo "address=/$domain/::/tag:$tag" >> "$DNS_CONF"
        idx=$((idx+1))
    done

    /etc/init.d/dnsmasq reload >/dev/null 2>&1
    elog "dns blacklist rebuilt (timed_only=$timed_only)"
}

del_rule() {
    for ch in parentcontrol_time parentcontrol_protocol parentcontrol_weburl parentcontrol_weburl_nat; do
        $nft flush chain inet fw4 $ch 2>/dev/null
        $nft delete chain inet fw4 $ch 2>/dev/null
    done
}

# cron 调用入口：阻断时间链（开启时间段）
block_time() {
    # 仅刷新时间链，避免影响协议/网址链
    $nft flush chain inet fw4 parentcontrol_time 2>/dev/null
    # 当前时间（分钟）与星期
    NOW_H=$(date +%H | sed 's/^0*//')
    NOW_M=$(date +%M | sed 's/^0*//')
    [ -z "$NOW_H" ] && NOW_H=0
    [ -z "$NOW_M" ] && NOW_M=0
    NOW_MIN=$((NOW_H*60 + NOW_M))
    NOW_WEEK=$(date +%w)
    [ "$NOW_WEEK" = "0" ] && NOW_WEEK=7   # 周日兼容
    TIME_FILTER=1
    set_rules time 0   # 常驻
    set_rules time 1   # 定时且在当前时间窗口内
    TIME_FILTER=0
    # 清理已有连接和可能的 flowtable 加速，避免旧连接继续穿透
    if command -v conntrack >/dev/null 2>&1; then
        conntrack -F 2>/dev/null
    fi
    $nft reset flowtable inet fw4 2>/dev/null
    touch "$LOCK"
    elog "cron: block_time executed (time rules active)"
}

# cron 调用入口：解除阻断（关闭时间段）
unblock_time() {
    $nft flush chain inet fw4 parentcontrol_time 2>/dev/null
    set_rules time 0   # 仅恢复常驻规则
    elog "cron: unblock_time executed (timed time rules removed, always rules kept)"
    # 注意：不删除 LOCK，让整体服务仍视为“启动”状态
}

block_protocol() {
    $nft flush chain inet fw4 parentcontrol_protocol 2>/dev/null
    set_rules protocol 0
    set_rules protocol 1
    elog "cron: block_protocol executed (timed protocol rules active)"
}

unblock_protocol() {
    $nft flush chain inet fw4 parentcontrol_protocol 2>/dev/null
    set_rules protocol 0
    elog "cron: unblock_protocol executed (timed protocol rules removed, always rules kept)"
}

block_weburl() {
    $nft flush chain inet fw4 parentcontrol_weburl 2>/dev/null
    $nft flush chain inet fw4 parentcontrol_weburl_nat 2>/dev/null
    NOW_H=$(date +%H | sed 's/^0*//'); NOW_M=$(date +%M | sed 's/^0*//')
    [ -z "$NOW_H" ] && NOW_H=0; [ -z "$NOW_M" ] && NOW_M=0
    NOW_MIN=$((NOW_H*60 + NOW_M))
    NOW_WEEK=$(date +%w); [ "$NOW_WEEK" = "0" ] && NOW_WEEK=7
    TIME_FILTER=1
    set_rules weburl 0
    set_rules weburl 1   # 仅加载当前时间窗口内的定时条目
    build_dns_blacklist ""   # 按 MAC 生成 tag 限定的 address 规则
    TIME_FILTER=0
    # 刷新连接，避免已建立 DNS/DoH 连接穿透
    if command -v conntrack >/dev/null 2>&1; then
        conntrack -F 2>/dev/null
    fi
    elog "cron: block_weburl executed (timed weburl rules active)"
}

unblock_weburl() {
    $nft flush chain inet fw4 parentcontrol_weburl 2>/dev/null
    $nft flush chain inet fw4 parentcontrol_weburl_nat 2>/dev/null
    NOW_H=$(date +%H | sed 's/^0*//'); NOW_M=$(date +%M | sed 's/^0*//')
    [ -z "$NOW_H" ] && NOW_H=0; [ -z "$NOW_M" ] && NOW_M=0
    NOW_MIN=$((NOW_H*60 + NOW_M))
    NOW_WEEK=$(date +%w); [ "$NOW_WEEK" = "0" ] && NOW_WEEK=7
    TIME_FILTER=1
    set_rules weburl 0
    set_rules weburl 1   # 仅加载当前时间窗口内的定时条目
    build_dns_blacklist 0
    TIME_FILTER=0
    elog "cron: unblock_weburl executed (timed weburl rules removed, always rules kept)"
}

start() {
    # 只要 enabled=1，就进入“启动”状态（即使无规则）
    enabled=$(uci -q get $CONFIG.@basic[0].enabled)
    [ "x$enabled" = "x1" ] || exit 1

    # 建日志目录与 lock
    mkdir -p "$LOG_DIR"
    touch "$LOCK"
    elog "parentcontrol service start requested"

    # 清理旧链，按配置加载“常驻”规则；定时规则交给 cron 触发 block_*
    del_rule

    # 统计启用条目
    allsum=$(grep -c "enable='1'" /etc/config/$CONFIG)
    if [ "$allsum" -gt 0 ]; then
        set_rules protocol 0
        set_rules weburl 0
        set_rules time 0
        build_dns_blacklist 0
        elog "rules loaded according to config (enabled sections=$allsum; timed rules wait cron)"
    else
        elog "no enabled sections; started without active rules"
    fi

    # 重建 cron（根据 time sections）
    rebuild_cron

    clean_log
}

stop() {
    # 彻底关闭：删链、删 lock、删 cron
    del_rule
    rm -f "$LOCK" 2>/dev/null
    sed -i '/# parentcontrol$/d;/\/etc\/init\.d\/parentcontrol \(\<block_time\>\|\<unblock_time\>\|\<block_protocol\>\|\<unblock_protocol\>\|\<block_weburl\>\|\<unblock_weburl\>\)/d' /etc/crontabs/root 2>/dev/null
    rm -f "$DNS_CONF" 2>/dev/null
    /etc/init.d/dnsmasq reload >/dev/null 2>&1
    /etc/init.d/cron restart
    elog "parentcontrol service stopped (rules/chains/cron cleared)"
}

reload() {
    enabled=$(uci get parentcontrol.@basic[0].enabled 2>/dev/null)

    if [ "$enabled" = "1" ]; then
        start
        # 立即激活规则（含定时条目），避免需手动 block_*
        block_protocol
        block_weburl
        block_time
    else
        stop
    fi
}

status() {
    enabled=$(uci -q get $CONFIG.@basic[0].enabled)

    # 1. 配置未启用 → 一律 inactive
    if [ "$enabled" != "1" ]; then
        echo "parentcontrol is inactive (disabled in config)"
        return 1
    fi

    # 2. 没有 lock → inactive
    if [ ! -f "$LOCK" ]; then
        echo "parentcontrol is inactive"
        return 1
    fi

    # 3. 有链 → running
    if $nft list chain inet fw4 parentcontrol_time >/dev/null 2>&1 || \
       $nft list chain inet fw4 parentcontrol_protocol >/dev/null 2>&1 || \
       $nft list chain inet fw4 parentcontrol_weburl >/dev/null 2>&1; then
        echo "parentcontrol is running (rules active)"
        return 0
    fi

    # 4. 无链但有 cron → running (等待触发)
    if grep -q "# parentcontrol" /etc/crontabs/root 2>/dev/null; then
        echo "parentcontrol is running (cron scheduled)"
        return 0
    fi

    # 5. 有 lock，但无链、无 cron → started(no rules)
    echo "parentcontrol is started (no rules)"
    return 0
}
